# linux/flowsuite.service
[Unit]
Description=FlowSuite local voice engine
After=graphical-session.target

[Service]
Type=simple
WorkingDirectory=__APPDIR__
ExecStart=__PYTHON3__ __APPDIR__/flowsuite.py
Restart=on-failure

[Install]
WantedBy=default.target
